# setup_files/autostart/pool-controller.desktop
[Desktop Entry]
Name=Pool Controller
Type=Application
Comment=Pool Controller Program
Exec=sh /home/user/Desktop/launcher.sh